# systemd/.config/systemd/user/helix-theme-sync.path
[Unit]
Description=Watch COSMIC theme changes

[Path]
PathModified=%h/.config/cosmic/com.system76.CosmicTheme.Mode/v1/is_dark

[Install]
WantedBy=default.target
